# infrastructure/svc-cache-registry/files/s6-rc.d/registry/run
#!/command/with-contenv sh
# Launch registry:2 with the same default config the upstream image ships
# (overridden via REGISTRY_* env vars on the compose service, e.g.
# REGISTRY_PROXY_REMOTEURL, REGISTRY_STORAGE_DELETE_ENABLED).
set -e
exec /entrypoint.sh /etc/docker/registry/config.yml
