# infrastructure/svc-cache-package/files/s6-rc.d/apt-cacher-ng/run
#!/command/with-contenv sh
set -e
exec s6-setuidgid apt-cacher-ng \
  /usr/sbin/apt-cacher-ng -c /etc/apt-cacher-ng foreground=1
